# examples/3-systemd-csv-rotation/cpu-thermals.service
[Unit]
Description=Continuous CPU thermal capture (cpu_thermals --csv -)
Documentation=https://github.com/cnh/cpu_thermals
After=network.target

# Cap restart attempts so a permanently broken config (wrong binary
# path, missing kernel module, etc.) fails the unit instead of looping
# forever and flooding the journal. Five retries in 60s is enough to
# absorb transient setup hiccups (e.g. a slow first sensor read on
# boot) without becoming noise.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple

# `--csv -` writes long-format CSV to stdout; systemd appends it to the log
# file via StandardOutput=append:. The TUI is auto-suppressed because
# stdout is the CSV.
ExecStart=/usr/local/bin/cpu-thermals --csv -

StandardOutput=append:/var/log/cpu_thermals/cpu_thermals.csv

# Banner ("recording CSV to stdout") and the Ctrl-C summary lines that
# cpu_thermals writes to stderr land in the journal, where they don't
# contaminate the CSV but are still visible via `journalctl -u cpu-thermals`.
StandardError=journal

# Restart only on actual failures (not on clean stop, e.g. when
# logrotate's postrotate script restarts us deliberately). Combined with
# the StartLimit* settings above, this fails fast on bad configs while
# still self-healing through transient errors.
Restart=on-failure
RestartSec=5

# Run as root for simplest sensors access. For hardening, swap in a
# dedicated user/group with read access to /sys/class/hwmon/* and write
# access to /var/log/cpu_thermals/. See README.md.
# User=cpu_thermals
# Group=cpu_thermals

# Optional: keep monitoring out of the way of real workloads.
# Nice=10
# IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
